# gentoo/init.d/sync-kernel-modules
#!/sbin/runscript
# Copyright (c) 2014 Online Labs Cloud Team <cloud-team@labs.online.net>
# Released under the 2-clause BSD license.

descriptionr="Synchronizes kernel modules to local filesystem"

depend()
{
	after root
	use labs.online.net.modules
}

start()
{
	# Should not fail if kernel do not have module
	# support compiled in ...
	[ ! -f /proc/modules ] && return 0

	local ret DIR KVERSION RAM_KMODULES
	DIR=/lib/modules
	KVERSION=$(uname -r)
	RAM_KMODULES="/run/modules"

	if grep -q "${DIR}" /proc/mounts
	then
		umount "${DIR}"
	fi

	if [ ! -d "${DIR}/${KVERSION}" -a -d "${RAM_KMODULES}/${KVERSION}" ]
	then
		ebegin "Writing modules to filesystem"
		mkdir -p "${DIR}/${KVERSION}"
		cp -a "${RAM_KMODULES}/${KVERSION}" "${DIR}"
		ret=$?
		if [ $ret -eq 0 ]
		then
			rm -rf "${RAM_KMODULES}"
		else
			rm -rf "${DIR}/${KVERSION}"
			mount -o ro,bind "${RAM_KMODULES}" "${DIR}"
		fi
		eend $ret
	fi
}
